# etc/udev/rules.d/99-input.rules
# Detecting bluetooth clients

SUBSYSTEM="input", GROUP="input", MODE="0660"
KERNEL=="input[0-9]*", RUN+="/home/paco/bluetooth_add.sh"
